# ESPNodes_frontend/bin/espnodes_daemon.service
[Unit]
Description=ESPNodes logging and timer daemon

[Service]
ExecStart=/usr/bin/php /var/www/html/espnodes/bin/daemon.php

[Install]
WantedBy=multi-user.target

